# public/configs/2-p2prun.conf
[program:P2P-Run]
command=bash -c 'p2pclient --login $MAIL_ADDRESS'
autostart=true
autorestart=true
priority=994
